# monitor_test.service
[Unit]
Description=Monitor Test Process
After=network.target

[Service]
ExecStart=/bin/bash /usr/local/bin/monitor_test.sh
Restart=always
User=root

[Install]
WantedBy=multi-user.target
